# copy/~/.config/autostart/clipboard-history.desktop
[Desktop Entry]
Type=Application
Name=Clipboard History
Exec=/home/hppc/Documents/copy/clipboard_logger.go
Hidden=false 




